# priod.service.in
[Unit]
Description=priod

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/priod @CMAKE_INSTALL_FULL_SYSCONFDIR@/priod.conf
Restart=on-abort

[Install]
WantedBy=default.target
